# 1.crawling/news_scrapy_crawling/crontab
0 0 * * * cd /app && scrapy crawl naver_news >> /var/log/cron.log 2>&1
